# init/twitter-stream.conf
start on runlevel [2345]

respawn
respawn limit 5 60

kill signal SIGINT

setuid vagrant
chdir /vagrant/twitter
exec php daemon.php